# lightswitch.service
[Unit]
Description=lightswitch

[Service]
ExecStart=/bin/lightswitch
